# config/sunzi/files/dtgcaa-upstart.conf
description "dtgcaa"

start on runlevel [2345]
stop on runlevel [016]

console log

exec start-stop-daemon --start -c dtgcaa --exec /home/dtgcaa/start.sh
